# lab7/autostart-docker.service
[Unit]
Description= Mon service autostarted
After=docker.service
Requires=docker.service
